# devops/workers/supervisor.conf
# https://github.com/Supervisor/supervisor/blob/main/supervisor/skel/sample.conf
[supervisord]
logfile=/tmp/supervisord.log ; main log file; default $CWD/supervisord.log
logfile_maxbytes=50MB        ; max main logfile bytes b4 rotation; default 50MB
logfile_backups=10           ; # of main logfile backups; 0 means none, default 10
loglevel=info                ; log level; default info; others: debug,warn,trace
pidfile=/tmp/supervisord.pid ; supervisord pidfile; default supervisord.pid
nodaemon=false               ; start in foreground if true; default false
silent=false                 ; no logs to stdout if true; default false
minfds=1024                  ; min. avail startup file descriptors; default 1024
minprocs=200                 ; min. avail process descriptors;default 200

[unix_http_server]
file=/tmp/supervisor.sock   ; the path to the socket file

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock ; use a unix:// URL  for a unix socket

###> Processos - Symfony Messenger ###
## "Fila" principal
[program:consume-async]
command=php bin/console messenger:consume async -vv --memory-limit=128M --time-limit=3600
directory=/var/www/html
stderr_logfile_maxbytes=0
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stdout_logfile=/dev/stdout
numprocs=1
startsecs=0
autostart=true
autorestart=true
startretries=10
process_name=%(program_name)s_%(process_num)02d

## "Fila" de mensagens com falha
[program:consume-failed]
command=php bin/console messenger:consume -vv failed --memory-limit=128M --time-limit=3600
directory=/var/www/html
stderr_logfile_maxbytes=0
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stdout_logfile=/dev/stdout
numprocs=1
startsecs=0
autostart=true
autorestart=true
startretries=10
process_name=%(program_name)s_%(process_num)02d

## CRON jobs - tarefas agendadas
[program:crons]
command=php bin/console messenger:consume -vv scheduler_cron --memory-limit=128M --time-limit=3600
directory=/var/www/html
stderr_logfile_maxbytes=0
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stdout_logfile=/dev/stdout
numprocs=1
startsecs=0
autostart=false
autorestart=true
startretries=10
process_name=%(program_name)s_%(process_num)02d
###< Processos - Symfony Messenger ###
